# data/cosmic-session.target
[Unit]
Description=Cosmic Session Target
Documentation=man:systemd.special(7)
BindsTo=graphical-session.target
Before=graphical-session.target
